# vlnd.service

[Unit]
Description=Virtual Local Network

[Service]
Type=simple
ExecStart=/usr/sbin/vlnd

[Install]
WantedBy=multi-user.target
